# test-utils/tob-spammer.service
[Unit]
Description=tob-spammer

[Service]
ExecStart=/home/ubuntu/pepc-boost-tester/pepc-tester

[Install]
WantedBy=multi-user.target 
