# test-utils/tob-spammer.service
[Unit]
Description=tob-spammer

[Service]
ExecStart=/home/ubuntu/pepc-boost-tester/tob-spammer

[Install]
WantedBy=multi-user.target 
